# systemd/systemd-ask-password-remdev.path
[Unit]
DefaultDependencies=no
Conflicts=shutdown.target
Before=paths.target shutdown.target

[Path]
DirectoryNotEmpty=/run/systemd/ask-password
MakeDirectory=yes
